# elasticsearch.conf
[program:elasticsearch]
command=/usr/share/elasticsearch/bin/elasticsearch
user=elasticsearch
autostart=true
autorestart=true
stderr_logfile=/var/log/elasticsearch/elasticsearch.log
