# kibana.conf
[program:kibana]
command=/opt/kibana/bin/kibana
user=kibana
autostart=true
autorestart=true
stderr_logfile=/var/log/kibana/kibana.stderr
